# docker-orientdb/supervisord.conf
[supervisord]
nodaemon=true

[program:orientdb]
command=/opt/orientdb/bin/server.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/orientdb.err.log
stdout_logfile=/var/log/orientdb.out.log
