# docker-jetty9-java8-download/supervisord.conf
[supervisord]
nodaemon=true

[program:jetty]
command=/opt/jetty/bin/jetty.sh run
autostart=true
autorestart=unexpected
